# medicine-bot.service
[Unit]
Description=Medicine Reminder Bot (ওষুধের রিমাইন্ডার)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/medicine-bot
EnvironmentFile=/home/ubuntu/medicine-bot/.env
ExecStart=/home/ubuntu/medicine-bot/venv/bin/gunicorn -w 1 -b 0.0.0.0:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
